# debian/creature-listener.service
[Unit]
Description=Creature Listener - Wake Word Conversational Interface
After=network-online.target sound.target
Wants=network-online.target

[Service]
EnvironmentFile=/etc/default/creature-listener
WorkingDirectory=/var/lib/creature-listener
ExecStart=/usr/bin/creature-listener \
    --config-path /etc/creature-listener.yaml
Restart=on-failure
RestartSec=5
# Run below creature-controller priority
Nice=5

[Install]
WantedBy=multi-user.target
